# debian/etc/cron.d/scripts
5 0 * * * root curl -X POST http://localhost:6101/cron
